# immich-dedup.timer
[Unit]
Description=Run Immich WhatsApp Duplicate Cleaner daily
Requires=immich-dedup.service

[Timer]
# Run daily at 5:00 AM
OnCalendar=*-*-* 05:00:00
# Run immediately if the system was down when the timer should have run
Persistent=true
# Add some randomization to avoid all systems running at exactly the same time
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
